# deploy/hermes/systemd/hermes-hourly.service
[Unit]
Description=Hermes hourly agentic discovery scout
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=hermes
Group=hermes
WorkingDirectory=/opt/vnem
EnvironmentFile=/etc/hermes/hermes.env
ExecStart=/usr/bin/env npm run hermes:hourly
Nice=5
