# deploy/hermes/systemd/hermes-hourly.timer
[Unit]
Description=Run Hermes hourly scout

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
